# systemd/proxy-exporter.service
[Unit]
Description=Proxy Exporter for Prometheus
After=network.target

[Service]
Type=simple
User=nobody
EnvironmentFile=/etc/proxy-exporter.env
ExecStart=/usr/bin/python3 /opt/monitoring/exporters/proxy_exporter.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
