# systemd/proxy-exporter.service
[Unit]
Description=Proxy Exporter for Prometheus
After=network.target

[Service]
Type=simple
User=nobody
EnvironmentFile=/etc/proxy-exporter.env
# Ubuntu 24.04 — externally-managed (PEP 668), aiohttp ставится в venv,
# а не системно. venv создаётся при деплое: /opt/monitoring/.venv
ExecStart=/opt/monitoring/.venv/bin/python /opt/monitoring/exporters/proxy_exporter.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
